# debian/cron-se.postinst
#!/bin/sh
set -e

# Fixup crontab binary, directory and files for new group 'crontab'.

if ! dpkg-statoverride --list /usr/bin/crontab > /dev/null ; then
    chgrp crontab /usr/bin/crontab
    chmod g+s /usr/bin/crontab
fi

if [ -x /etc/init.d/cron ]; then
	if [ -x /usr/sbin/invoke-rc.d ] ; then
	       invoke-rc.d cron restart
	else
	       /etc/init.d/cron restart
	fi
fi



#DEBHELPER#
