# debian/cron.upstart.in
# cron - regular background program processing daemon
#
# cron is a standard UNIX program that runs user-specified programs at
# periodic scheduled times

description	"regular background program processing daemon"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn

script
    for ENV_FILE in /etc/environment /etc/default/locale; do
        [ -r "$ENV_FILE" ] || continue
        [ -s "$ENV_FILE" ] || continue

        for var in LANG LANGUAGE LC_ALL LC_TYPE; do
            value=`egrep "^${var}=" "$ENV_FILE" | tail -n1 | cut -d= -f2`
            [ -n "$value" ] && eval export $var=$value
        done
    done

    # Force collation sequence since ASCII is expected in regexps
    export LC_COLLATE=C

    if [ -z "$TZ" -a -e /etc/timezone ]; then
        export TZ=`cat /etc/timezone`
    fi

    if [ -f /etc/default/cron ]; then
        . /etc/default/cron
    fi

    # These variables may be set in /ec/default/cron
    exec cron $LSBNAMES $EXTRA_OPTS
end script
